# images/docker/ubuntu-24.04/overlay/etc/systemd/system/fireactions.service
[Unit]
Description=Fireactions
Documentation=https://github.com/hostinger/fireactions
# Wait for cloud-init to complete - CA certificate and DNS config are set via cloud-init runcmd
# Note: After= is ignored if cloud-final.service doesn't exist (e.g., no cloud-init installed)
After=network.target cloud-final.service
SuccessAction=reboot

[Service]
Type=simple
User=root
ExecStart=/usr/bin/fireactions runner --log-level=info
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
