# images/docker/ubuntu-24.04/overlay/etc/systemd/system/fireactions.service
[Unit]
Description=Fireactions
Documentation=https://github.com/hostinger/fireactions
# Wait for cloud-init config stage - CA cert and DNS are set via runcmd in cloud-config.service
# Note: Can't use cloud-final.service - it runs AFTER multi-user.target, causing ordering cycle
# After= is ignored if service doesn't exist (e.g., no cloud-init installed)
After=network.target cloud-config.service
SuccessAction=reboot

[Service]
Type=simple
User=root
ExecStart=/usr/bin/fireactions runner --log-level=info
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
