# images/docker/ubuntu-24.04/overlay/etc/systemd/system/fireactions.service
[Unit]
Description=Fireactions
Documentation=https://github.com/hostinger/fireactions
# Wait for cloud-init.service which runs ca_certs module
# ca_certs handles both writing the cert AND running update-ca-certificates
# After= is ignored if service doesn't exist (e.g., no cloud-init installed)
After=network.target cloud-init.service
SuccessAction=reboot

[Service]
Type=simple
User=root
ExecStart=/usr/bin/fireactions runner --log-level=info
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
