# images/docker/ubuntu-24.04/overlay/gitlab/etc/systemd/system/fireglab.service
[Unit]
Description=Fireglab GitLab CI Runner Agent
Documentation=https://github.com/thpham/nixos-fireactions
# Wait for cloud-config.service which runs runcmd
# Timeline:
# 1. cloud-init.service (init stage): ca_certs writes cert + creates symlinks
# 2. cloud-config.service (config stage): runcmd appends cert to bundle file
#    (workaround for Ubuntu 24.04 bug where ca_certs doesn't update bundle)
# 3. Fireglab starts with CA bundle containing our cert
# After= is ignored if service doesn't exist (e.g., no cloud-init installed)
After=network.target cloud-config.service docker.service
Requires=docker.service
SuccessAction=reboot

[Service]
Type=simple
User=root
# The runner command:
# 1. Fetches metadata from MMDS (glrt-* token, gitlab URL, tags)
# 2. Registers gitlab-runner with GitLab using the token
# 3. Starts gitlab-runner daemon in single-job mode
# 4. Exits when job completes (triggering VM reboot via SuccessAction)
ExecStart=/usr/bin/fireglab runner --log-level=info --single-job
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
